# root/etc/s6-overlay/s6-rc.d/init-kasmvnc/up
# This file doesn't do anything, it's just the beginning of the kasmvnc init process