# root/etc/s6-overlay/s6-rc.d/svc-de/run
#!/usr/bin/with-contenv bash

cd $HOME
exec s6-setuidgid abc \
  /bin/bash /defaults/startwm.sh
