# root/etc/s6-overlay/s6-rc.d/svc-kclient/run
#!/usr/bin/with-contenv bash

# NodeJS wrapper
cd /kclient
exec s6-setuidgid abc \
  node index.js
