# root/etc/s6-overlay/s6-rc.d/svc-kasmvnc/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc \
  /usr/local/bin/Xvnc $DISPLAY \
    -disableBasicAuth \
    -SecurityTypes None \
    -AlwaysShared \
    -http-header Cross-Origin-Embedder-Policy=require-corp \
    -http-header Cross-Origin-Opener-Policy=same-origin \
    -geometry 1024x768 \
    -sslOnly 0 \
    -RectThreads 0 \
    -websocketPort 6901 \
    -interface 0.0.0.0 \
    -Log *:stdout:10
