# root-armhf/etc/s6-overlay/s6-rc.d/svc-splash/run
#! /bin/bash

echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'
echo 'THIS IMAGE DOES NOT WORK ON ARMHF'

s6-notifyoncheck -d -n 300 -w 1000 -c "exit 0" sleep infinity
